# centos7-develop/s6/php-fpm/run
#!/bin/sh

set -e

PHP_RUN_PID='/var/run/php-fpm/php-fpm.pid'

# Cleanup stale pid on restart
rm -f "${PHP_RUN_PID}"

# Run php fpm in foreground
exec /usr/sbin/php-fpm --nodaemonize
