# centos7-develop/s6/httpd/run
#!/bin/sh

set -e

APACHE_RUN_PID='/var/run/httpd/httpd.pid'

# Cleanup stale pid on restart
rm -f "${APACHE_RUN_PID}"

# Run apache httpd in foreground
exec /usr/sbin/httpd -DFOREGROUND
